# contrib/init/dogecoind.conf
description "Nicoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env NICOIND_BIN="/usr/bin/nicoind"
env NICOIND_USER="nicoin"
env NICOIND_GROUP="nicoin"
env NICOIND_PIDDIR="/var/run/nicoind"
# upstart can't handle variables constructed with other variables
env NICOIND_PIDFILE="/var/run/nicoind/nicoind.pid"
env NICOIND_CONFIGFILE="/etc/nicoin/nicoin.conf"
env NICOIND_DATADIR="/var/lib/nicoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # nicoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$NICOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run nicoind."
        echo "The setting must appear in $NICOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Nicoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$NICOIND_PIDDIR"
    chmod 0755 "$NICOIND_PIDDIR"
    chown $NICOIND_USER:$NICOIND_GROUP "$NICOIND_PIDDIR"
    chown $NICOIND_USER:$NICOIND_GROUP "$NICOIND_CONFIGFILE"
    chmod 0660 "$NICOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$NICOIND_PIDFILE" \
    --chuid $NICOIND_USER:$NICOIND_GROUP \
    --exec "$NICOIND_BIN" \
    -- \
    -pid="$NICOIND_PIDFILE" \
    -conf="$NICOIND_CONFIGFILE" \
    -datadir="$NICOIND_DATADIR" \
    -disablewallet \
    -daemon

